# deployment/fillnode-worker.target
[Unit]
Description=Fillnode Background Worker
Wants=fillnode-worker.1.service

[Install]
WantedBy=multi-user.target